# etc/init.d/vdr
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/media-video/vdr/files/rc.vdr,v 1.23 2005/05/11 20:17:42 zzam Exp $

opts="${opts} watchdogrestart"

vdr_rcdir=/usr/lib/vdr/rcscript
source ${vdr_rcdir}/functions.sh

depend() {
	use lircd dvbsplash
	need net
	after checkroot
}

start() {
	einfo "Preparing start of vdr:"
	getvdrversion
	init_params
	init_daemonctrl_params

	[[ -z "${VDR_BIN}" ]] && VDR_BIN=/usr/bin/vdr
	add_daemonctrl_param --start --chdir ~vdr --chuid vdr --exec ${VDR_BIN}

	load_addons_prefixed pre-start || return 1

	for PLUGIN in ${PLUGINS}; do
		load_plugin ${PLUGIN}
	done

	if [[ "${waitconditions}" ]]; then
		ebegin "  Prerequisits (devices nodes etc.)" 
		waitfor 10 wait_for_multiple_condition
		exitcode="$?"
		exitmsg="could not start vdr"
		case "$exitcode" in
			1)
				exitmsg="${exitmsg}: ${condition_msg}"
				;;
			2)
				exitmsg="${exitmsg}: Timeout, ${condition_msg}"
				;;
		esac
		eend "$exitcode" "${exitmsg}" || return
	fi

	ebegin "Starting vdr-${vdrversion}"

	#
	# finally start vdr
	# 
	/usr/bin/dvbsplash stop 2>/dev/null >/dev/null
	unset MAIL
	export LD_ASSUME_KERNEL=2.4.1
	einfo_level2 setting LD_ASSUME_KERNEL=2.4.1

	vdr_home=/var/vdr

	export LOGNAME=vdr
	export USER=vdr
	export HOME=${vdr_home}

	cd ${vdr_home}

	einfo_level2 start-stop-daemon "${daemonctrl_opts[@]}" -- "${vdr_opts[@]}"

	if [[ -z "${TERMINAL}" ]]; then
		start-stop-daemon "${daemonctrl_opts[@]}" \
			-- --daemon "${vdr_opts[@]}"
		exitcode=$?
	else
		TERMINAL=${TERMINAL/\/dev\/tty/}
		TERMINAL_DEVICE=/dev/tty${TERMINAL}
		{
			clear
			einfo "Starting vdr"
		} >${TERMINAL_DEVICE}

		# strange commandline, to be fixed in next versions
		openvt -c ${TERMINAL} -- \
			start-stop-daemon "${daemonctrl_opts[@]}" \
			-- "${vdr_opts[@]}"
		exitcode=$?
	fi
	eend $exitcode "Failed to start vdr."

	if [[ "${exitcode}" == "0" ]]; then
		load_addons_prefixed post-start || exitcode=1
	fi
	return $exitcode
}

stop() {
	load_addons_prefixed pre-stop

	getvdrversion
	ebegin "Stopping vdr-${vdrversion}"

	# Use --name here to allow us to kill vdr even after a new emerge
	start-stop-daemon --stop --quiet --retry 15 --name vdr
	eend $? "Failed to stop vdr."

	load_addons_prefixed post-stop
}

watchdogrestart() {

	export WATCHDOG_RESTART="1"
	/etc/init.d/vdr pause
	sleep 2
	test_vdr_process && /bin/killall -9 vdr
	sleep 1

	/etc/init.d/vdr zap

	load_addons_prefixed watchdog-restart

	/etc/init.d/vdr start
	unset WATCHDOG_RESTART
}
