# etc/init.d/vdr
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License v2
# $Id$

opts="${opts} watchdogrestart"

. /etc/conf.d/vdr.watchdogd

load_functions() {
	cd /
	. /usr/share/vdr/inc/functions.sh
	include rc-functions
}

#
# Used to log error-messages in startscript to show them on
# OSD later when choosing apropriate point in commands.
#

vdr_log()
{
	echo "$@" >> ${VDR_LOG_FILE}
}

depend() {
	need net
	[ "${IR_CTRL}" = "lirc" ] && need lircd
	use lircd coldplug
	after checkroot
}

start_vdr() {
	load_functions
	local exitcode

	VDR_LOG_FILE=/var/vdr/tmp/vdr-start-log
	rm -f "${VDR_LOG_FILE}"
	> "${VDR_LOG_FILE}"

	einfo "Preparing start of vdr:"
	init_params
	add_daemonctrl_param --start --chdir ~vdr --exec ${VDR_BIN}
	
	load_addons_prefixed pre-start || return 1

	ebegin "Starting ${VDRNAME}"

	#
	# finally start vdr
	# 
	unset MAIL

	vdr_home=/var/vdr

	export LOGNAME=vdr
	export USER=vdr
	export HOME=${vdr_home}

	cd ${vdr_home}

	debug_msg "  CMDLINE:" start-stop-daemon "${daemonctrl_opts}" "--" "${vdr_opts}"

	if [ -z "${TERMINAL}" ]; then
		eval start-stop-daemon ${daemonctrl_opts} \
			-- --daemon ${vdr_opts}
		exitcode=$?
	else
		TERMINAL=${TERMINAL##/dev/tty}
		TERMINAL_DEVICE=/dev/tty${TERMINAL}
		{
			clear
			einfo "Starting vdr"
		} >${TERMINAL_DEVICE}

		# strange commandline, to be fixed in next versions
		openvt_opts=""
		if [ "${SWITCH_TO_TERMINAL}" = "yes" ]; then
			openvt_opts="${openvt_opts} -s"
		fi
		eval openvt -c ${TERMINAL} ${openvt_opts} -- \
			start-stop-daemon ${daemonctrl_opts} \
			-- ${vdr_opts}
		exitcode=$?
	fi
	eend $exitcode "Failed to start vdr."

	if [ "${exitcode}" = "0" ]; then
		load_addons_prefixed post-start
		if [ "$?" != "0" ]; then
			exitcode=1
		fi
	fi

	# wenn nachrichten vorhanden sind
	if [ -s "${VDR_LOG_FILE}" ]; then
		/usr/share/vdr/bin/vdr-bg.sh svdrpsend.pl mesg "Errors: Go to Commands/View VDR Start Log"
	fi

	return $exitcode
}

stop_vdr() {
	load_functions
	load_addons_prefixed pre-stop

	ebegin "Stopping ${VDRNAME}"

	# Use --name here to allow us to kill vdr even after a new emerge
	start-stop-daemon --stop --quiet --retry 15 --exec ${VDR_BIN} 
	exitcode=$?
	eend $exitcode "Failed to stop vdr."

	load_addons_prefixed post-stop
	return $exitcode
}

start_watchdog() {
	[ "${ENABLE_EXTERNAL_WATCHDOG:-yes}" = "yes" ] || return 0
	ebegin "Starting vdr watchdog"
	start-stop-daemon \
		--start \
		--background \
		--make-pidfile \
		--pidfile /var/run/vdrwatchdog.pid \
		--exec /usr/sbin/vdr-watchdogd \
		--name vdr-watchdogd
	eend $? "failed starting vdr watchdog"
	return 0
}


stop_watchdog() {
	if [ "${ENABLE_EXTERNAL_WATCHDOG:-yes}" = "yes" ]; then
		ebegin "Stopping vdr watchdog"
		start-stop-daemon --stop --pidfile /var/run/vdrwatchdog.pid --name vdr-watchdogd
		eend $? "failed stopping watchdog"
	else
		# Also stop watchdog if conf was changed to disabled while it was running
		start-stop-daemon --stop --pidfile /var/run/vdrwatchdog.pid --name vdr-watchdogd --quiet
	fi
	return 0
}

start() {
	start_vdr || return 1
	start_watchdog || return 2
}

stop() {
	stop_watchdog || return 2
	stop_vdr
}

# gets called by watchdog to restart vdr
# and possibly reload modules
watchdogrestart() {
	load_functions

	( stop_vdr )
	if test_vdr_process; then
		sleep 2
		test_vdr_process && killall -9 vdr
	fi

	load_addons_prefixed watchdog-restart

	( start_vdr )
}
