# etc/init.d/vdr-x11
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Start minimal X server for VDR"

. /etc/conf.d/vdr

command="/usr/bin/X"
command_args="${X11_SERVER:=\:0 -br -nocursor -nolisten tcp -quiet}"
start_stop_daemon_args="--make-pidfile --background"
pidfile="/var/run/vdr-x11.pid"

depend() {
	need localmount
	after acpid
	before vdr
}
