# etc/init.d/vdr-frontend
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

command="/usr/bin/vdr-sxfe"
command_args="--silent --fullscreen --syslog --reconnect xvdr://localhost &>/dev/null"
pidfile="/var/run/vdr-frontend.pid"
start_stop_daemon_args="--user vdr --group vdr --make-pidfile --background"

depend() {
	need vdr-x11
	after vdr-x11
}
