# etc/init.d/wakeup-reboot-halt
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need localmount
}

start() {
	# check weather this time this is a reboot and I should shutdown
	
	local STATEFILE=/var/vdr/shutdown-data/special_script_should_shutdown
	if [[ -e ${STATEFILE} ]]; then
		einfo "This is a reboot to write the wakeup-time into the bios."
		einfo "Shutting down now."
		rm -f ${STATEFILE}
		init 0
	fi
}

stop() {
	:
}

restart() {
	:
}
